# infra/src/ec2/minecraft-dns.service
[Unit]
Description=Update Cloudflare DNS for Minecraft server on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-dns.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target