# infra/src/ec2/minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
# Treat the screen process as the long-running process; don't expect a fork/PID file
Type=simple
User=minecraft

WorkingDirectory=/opt/minecraft/server
PermissionsStartOnly=true

# Mark /opt/setup as safe (git ownership check workaround)
ExecStartPre=/usr/bin/git config --global --add safe.directory /opt/setup

# Update the repo to get the latest config
ExecStartPre=/usr/bin/git -C /opt/setup pull

# Sync only your repo files (won’t delete world/)
ExecStartPre=/usr/bin/rsync -a /opt/setup/config/ /opt/minecraft/server/

# Ensure correct ownership
ExecStartPre=/bin/chown -R minecraft:minecraft /opt/minecraft/server/

# Start inside a detached screen session
ExecStart=/usr/bin/screen -DmS mc-server /usr/bin/java -Xms3276M -Xmx3276M -jar paper.jar nogui

# Gracefully ask Minecraft to stop
ExecStop=/usr/bin/screen -S mc-server -p 0 -X stuff "stop$(printf \\r)"

# Then kill off the screen session itself
ExecStopPost=/usr/bin/screen -S mc-server -X quit || true

# Ensure systemd kills any remaining children of this unit
KillMode=control-group

# Give startup extra time; worlds can take a while to load
TimeoutStartSec=300
# Give Minecraft up to 60s to shut down cleanly
TimeoutStopSec=60

Restart=on-failure

# Security hardening
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/minecraft

[Install]
WantedBy=multi-user.target
